# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:rundeck]
command=/etc/init.d/rundeckd start
autostart=true
startsecs=0
